# packaging/linux/systemd/soteria.service
[Unit]
Description=Soteria Aegis Encrypted Security Platform
After=network.target local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/soteriad serve --bind 127.0.0.1 --port 7777
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/etc/soteria /var/lib/soteria
SupplementaryGroups=tss
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soteria

[Install]
WantedBy=multi-user.target
